# build_files/files/car-edge-setup-wizard.service
[Unit]
Description=Bazzite Car Edge - First Boot Setup Wizard
After=graphical-session.target
ConditionPathExists=!%h/.config/car-edge-setup-complete

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/car-edge-setup-wizard --auto
RemainAfterExit=yes

[Install]
WantedBy=default.target
